# .config/supervisord.conf
[program:nginx]
command=nginx -g 'daemon off;'

[program:uwsgi]
command=/home/ubuntu/.local/bin/uwsgi --ini /srv/projects/healthcare/.config/uwsgi.ini